# knapsack-service/knapsack-client/supervisord.conf
[supervisord]
nodaemon=true

[program:knapsack-client]
command=java -jar -Xms4096m -Xmx4096m /data/knapsack-client.jar server /data/configuration.yaml